# douyin-api.service
[Unit]
Description=Douyin API Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/app/scrapper_douyin/douyin_api
ExecStart=/app/scrapper_douyin/.venv/bin/python start.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
